# luci-app-singbox-ui/root/etc/uci-defaults/50_luci-singbox-ui
#!/bin/sh

# Remove old misnamed ACL file left by previous versions
rm -f /usr/share/rpcd/acl.d/luci-app-singb.json

# Clear LuCI caches so the new app appears immediately
rm -f /var/luci-indexcache* /tmp/luci-indexcache* 2>/dev/null
find /tmp -name 'luci-*cache*' -delete 2>/dev/null

# Restart rpcd to pick up the new ACL
[ -x /etc/init.d/rpcd ] && /etc/init.d/rpcd restart

# Restart uhttpd to serve updated static resources
[ -x /etc/init.d/uhttpd ] && /etc/init.d/uhttpd restart

exit 0
